# config/crontab
0 11 * * *  cd /home/teflonted/marina.anachromystic.com/current && nice -n 19 rake videos:latest RAILS_ENV=production
